# supervisor/app.conf
[program:waveforms]
environment=DISPLAY=:1
command=waveforms